# src/modules/mistos/filesystem/root/etc/systemd/system/mistos_update_lightdm.service
[Unit]
Description=MistOS: Update LightDM config on first boot.

[Service]
ExecStart=/etc/mistos/scripts/update-lightdm.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=lightdm.service